# roles/packages/files/packages.cron
#!/bin/bash

cd /srv/docker/ansible/
ansible-playbook local.yml --tags=packages -c local
